# ddns_update.cron
#!/bin/bash

export PATH=$PATH:/sbin:/usr/sbin:/usr/local/bin

/opt/ddns_update/ddns_update.sh

